# deploy/systemd/giga-frontend.service
[Unit]
Description=GIGA Interior Next.js Frontend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/app/frontend
EnvironmentFile=/var/app/frontend/.env.frontend
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

